# deploy/warehouse.conf
[program:warehouse]
command=/usr/bin/python3 -m uvicorn main:app --host 127.0.0.1 --port 8000
directory=/var/www/warehouse/backend
autostart=true
autorestart=true
stderr_logfile=/var/log/warehouse.err.log
stdout_logfile=/var/log/warehouse.out.log
user=www-data
environment=PYTHONUNBUFFERED=1
